# docker/php/etc/supervisor/conf.d/messenger-worker.ini
;/etc/supervisor/conf.d/messenger-worker.conf
[program:messenger-consume]
command=php /var/www/html/apps/MusicLabel/api/bin/console messenger:consume async --time-limit=3600
user=root
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d